# deploy/algotrading-token-alert.service
[Unit]
Description=AlgoTrading Token Expiry Alert (08:30-09:15 IST)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=algotrading
WorkingDirectory=/opt/algotrading/app
ExecStart=/opt/algotrading/app/venv/bin/python token_alert_service.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
